# tools/pritunl-cloud-redirect.socket
[Unit]
Description=Pritunl Cloud Redirect Server Socket

[Socket]
ListenStream=80
Accept=no

[Install]
WantedBy=sockets.target
