# tools/pritunl-cloud.service
[Unit]
Description=Pritunl Cloud Daemon

[Service]
LimitNOFILE=50000
ExecStart=/usr/bin/pritunl-cloud start

[Install]
WantedBy=multi-user.target
